# rvc-core.service.example
[Unit]
Description=ReVoiceChat Core Server

[Service]
Type=simple
User=administrateur
Group=administrateur

WorkingDirectory=[PATH TO YOUR DIR]
ExecStart=/usr/bin/java -Dquarkus.config.locations=server.properties -jar [PATH TO YOUR DIR]/app/target/quarkus-app/quarkus-run.jar

[Install]
WantedBy=multi-user.target
