# rvc-core.service.example
[Unit]
Description=ReVoiceChat Core Server

[Service]
Type=simple
User=administrateur
Group=administrateur

WorkingDirectory=[PATH TO YOUR DIR]
ExecStart=[PATH TO YOUR DIR]/mvnw spring-boot:run -pl "app" -Dspring-boot.run.profiles=linux,pg

[Install]
WantedBy=multi-user.target
